# packaging/htu21.service
[Unit]
Description=HTU21 JSON HTTP endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=htu21
Group=htu21
WorkingDirectory=/opt/htu21
ExecStart=/opt/htu21/htu21
Environment=HTU21_PORT=273
Environment=DOTNET_BUNDLE_EXTRACT_BASE_DIR=/opt/htu21/.net
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target